# deploy/ddns-server.service
[Unit]
Description=DDNS Server Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=duacodie
Group=duacodie
WorkingDirectory=/opt/duacodie
EnvironmentFile=/opt/duacodie/.env
ExecStart=/opt/duacodie/bin/ddns start
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
